# supervisord/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
redirect_stderr=true
